# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f
stdout_logfile=/var/log/cron.log
stderr_logfile=/var/log/cron_error.log

[program:gunicorn]
command=gunicorn -k uvicorn.workers.UvicornWorker -w 5 -b 0.0.0.0:8000 --timeout 3600 --chdir ./src server:app
stdout_logfile=/var/log/gunicorn.log
stderr_logfile=/var/log/gunicorn_error.log
